# systemd/tas2781-firmware-reload.service
[Unit]
Description=TAS2781 amplifier boot reset (PCI D3cold)
# The BIOS initialises the TAS2781 amps in a state where direct I2C
# register programming via i2cset (2pa-byps.sh) may not take full effect.
# Removing the I2C controller PCI device (0000:00:15.2) causes ACPI to gate
# its power rail (D3cold), resetting all downstream I2C devices including the
# two TAS2781 amps to factory state. A subsequent PCI rescan re-enumerates the
# controller so 2pa-byps.sh can program the amps cleanly when WirePlumber
# starts (via tas2781-amp-fix.service).
#
# Note: snd_hda_scodec_tas2781_i2c is blacklisted via
# /etc/modprobe.d/tas2781-blacklist.conf - the modprobe step below is a
# safe no-op kept for belt-and-suspenders safety.
#
# Must complete before display-manager starts: yanking the PCI device
# while WirePlumber holds an ALSA handle causes a SEGV in libasound
# (snd_hctl_elem_get_interface on element re-add). Running before
# display-manager ensures PipeWire has not yet opened the device.
After=sound.target
Before=display-manager.service

[Service]
Type=oneshot
# Give the initial driver probe time to complete before removing the device.
ExecStartPre=/bin/sleep 2
# Remove the I2C controller PCI device. This tears down i2c-4 and all
# child devices (i2c-TIAS2781:00), then ACPI power-gates the controller.
ExecStart=/bin/bash -c 'echo 1 > /sys/bus/pci/devices/0000:00:15.2/remove'
# Allow the hardware power-off to settle so the amps reset to factory state.
ExecStart=/bin/sleep 2
# Re-enumerate the I2C controller. udevd will rebind i2c_designware and
# recreate i2c-4.
ExecStart=/bin/bash -c 'echo 1 > /sys/bus/pci/rescan'
# Blacklisted - this is a no-op, but kept so the service degrades gracefully
# if the blacklist is ever removed.
ExecStart=/bin/bash -c 'modprobe snd_hda_scodec_tas2781_i2c || true'
# Allow i2c_designware to finish binding before display-manager (and
# PipeWire/WirePlumber) start.
ExecStart=/bin/sleep 5
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
